# deploy/ec2/powerclaw-backup.timer
[Unit]
Description=Run PowerClaw state backup daily

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=20m

[Install]
WantedBy=timers.target
